# backend/nemnovo-backend.service
# systemd unit для Django (gunicorn)
# Gunicorn слушает 127.0.0.1:8000 — nginx проксирует на него.
#
# Установка:
#   1. Замените /home/user/nemnovo на реальный путь к backend
#   2. sudo cp nemnovo-backend.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable nemnovo-backend
#   5. sudo systemctl start nemnovo-backend

[Unit]
Description=Nemnovo Tour Django Backend (Gunicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/user/nemnovo/backend
Environment="PATH=/home/user/nemnovo/backend/.venv/bin"
ExecStart=/home/user/nemnovo/backend/.venv/bin/gunicorn config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 2 \
    --threads 2 \
    --timeout 60 \
    --access-logfile - \
    --error-logfile -
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
